# MyService.service
[Unit]
Description=My New Service

[Service]
ExecStart=/root/src/start.sh

[Install]
WantedBy=default.target
