# MyService.service
[Unit]
Description=My New Service

[Service]
ExecStart=~/src/start.sh

[Install]
WantedBy=default.target
